# build/overlay/etc/cron.d/dongle-monitoring
# Signal quality every 5 minutes
*/5 * * * * root /usr/local/bin/signal-monitor.sh

# Connection watchdog every 3 minutes
*/3 * * * * root /usr/local/bin/connection-watchdog.sh

# Data usage every hour
0 * * * * root /usr/local/bin/data-usage.sh

# Sync clock daily (no RTC battery on the dongle)
0 3 * * * root date -s "$(curl -sI http://google.com 2>/dev/null | grep -i '^date:' | cut -d' ' -f2-)" 2>/dev/null || true
