# build/overlay/etc/systemd/system/clock-sync.service
[Unit]
Description=Sync system clock via HTTP Date header
After=modem-autoconnect.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'for i in 1 2 3 4 5; do DATE=$(wget -qO /dev/null --save-headers http://google.com 2>&1 | grep -i "^Date:" | cut -d" " -f2-); if [ -n "$DATE" ]; then date -s "$DATE" && echo "Clock synced: $(date)" && exit 0; fi; sleep 10; done; echo "Clock sync failed"'

[Install]
WantedBy=multi-user.target
