# build/overlay/etc/systemd/system/modem-autoconnect.service
[Unit]
Description=Auto-connect LTE modem and configure wwan0
After=ModemManager.service rmtfs.service
Wants=ModemManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/modem-autoconnect.sh
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
