# build/overlay/etc/systemd/system/dnsmasq.service.d/after-gadget.conf
[Unit]
After=usb-gadget.service
Requires=usb-gadget.service
